# aiden-bot.service
[Unit]
Description=Aiden Telegram Bot
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/aiden bot
Environment="PATH=/home/%i/aiden bot/venv/bin"
ExecStart=/home/%i/aiden bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=append:/home/%i/aiden bot/logs/bot.log
StandardError=append:/home/%i/aiden bot/logs/bot.error.log

# Лимиты ресурсов
LimitNOFILE=65536
Nice=10

[Install]
WantedBy=multi-user.target
